# packaging/deb/postinst
#!/usr/bin/env bash
# xdr-agent - Modular XDR endpoint security agent for Linux
# Copyright (C) 2026  Diego A. Guillen-Rosaperez
# SPDX-License-Identifier: AGPL-3.0
set -e

if [ ! -f /etc/xdr-agent/config.json ] && [ -f /usr/share/xdr-agent/config.default.json ]; then
  mkdir -p /etc/xdr-agent
  cp -f /usr/share/xdr-agent/config.default.json /etc/xdr-agent/config.json
  chmod 0644 /etc/xdr-agent/config.json
fi

if [ -x /usr/bin/xdr-agent ]; then
  mkdir -p /usr/share/bash-completion/completions || true
  if [ -f /usr/share/xdr-agent/bash_completion/xdr-agent ]; then
    cp -f /usr/share/xdr-agent/bash_completion/xdr-agent /usr/share/bash-completion/completions/xdr-agent || true
  fi
  chmod 0644 /usr/share/bash-completion/completions/xdr-agent || true
fi

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload || true
  systemctl preset xdr-agent.service || true
  systemctl disable xdr-agent.service || true
  systemctl stop xdr-agent.service || true
fi

exit 0
